# apps/bot/deploy/bnb-bot.service
[Unit]
Description=Bits&Bytes Discord Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/bits-bytes-bot
EnvironmentFile=/opt/bits-bytes-bot/.env
Environment=BUN_INSTALL=/home/ubuntu/.bun
ExecStart=/home/ubuntu/.bun/bin/bun run index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
