# deploy/api/bnb-api.service
[Unit]
Description=BNB API (FastAPI/Uvicorn)
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/bnb-api/apps/api
EnvironmentFile=/opt/bnb-api/.env
ExecStart=/opt/bnb-api/apps/api/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1 --limit-max-requests 10000
Restart=always
RestartSec=5

# Strict resource limits for 300MB VPS
MemoryMax=200M
MemoryHigh=150M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=bnb-api

[Install]
WantedBy=multi-user.target
